# debian/postinst
#!/bin/bash

update-rc.d watcherd defaults
service watcherd start
