# debian/prerm
#!/bin/bash

update-rc.d -f watcherd remove
service watcherd stop
